# recipes-bsp/hm-legacy-names/files/hm-legacy-names.service
[Unit]
Description=Legacy Network Names Service
After=systemd-udevd.service
Before=network.target

[Service]
ExecStart=/opt/hm/legacy-network-names.sh

[Install]
WantedBy=network.target
