# recipes-flexray-utils/fr-mcu/files/fr-mcu.rules
SUBSYSTEM=="net", ACTION=="add", KERNEL=="eth2",  RUN+="/bin/sh -c '/sbin/ifconfig eth2 10.42.42.1;' "
