# recipes-bsp/udev/udev-rules-imx/imx8mp-var-dart/99-hm-interface-rename.rules
# Trigger rename script when tmp_ interfaces are created
SUBSYSTEM=="net", KERNEL=="tmp_*", RUN+="/opt/hm/hm-rename-interface.sh %k"
