# recipes-bsp/hm-autostart/verdin-am62-hmm/autostart.timer
[Unit]
Description=Run autostart 30 seconds after boot

[Timer]
OnBootSec=30s
Unit=autostart.service

[Install]
WantedBy=timers.target