# recipes-bsp/hm-autostart/hm-autostart/autostart.service
[Unit]
Description=Run autostart script

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/hm/autostart.sh