# recipes-bsp/hm-autostart/hm-autostart/autostart.timer
[Unit]
Description=Run autostart 15 seconds after boot

[Timer]
OnBootSec=15s
Unit=autostart.service

[Install]
WantedBy=timers.target