# recipes-bsp/hm-interface-names/files/imx8mp-var-dart-hmx1/10-hm-interface-names.rules
SUBSYSTEM=="net", ACTION=="add" RUN+="systemctl start hm-interface-names"

